# systemd/bitlog.timer
[Unit]
Description=logging worker process A

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=timers.target
